# uci-defaults/luci-i18n-ukrainian
#!/bin/sh
uci batch <<-EOF
	set luci.languages.uk=Україна
        commit luci
EOF

