# uci-defaults/luci-i18n-english
#!/bin/sh
uci batch <<-EOF
	set luci.languages.en=English
        set wireless.wifi*.country='EN'
        commit luci
EOF

