# uci-defaults/luci-i18n-turkish
#!/bin/sh
uci batch <<-EOF
	set luci.languages.tr=Türkçe
        commit luci
EOF

