# usr/lib/opkg/info/conntrackd.list
/usr/sbin/conntrackd
/lib/upgrade/keep.d/conntrackd
/etc/init.d/conntrackd
