# usr/lib/opkg/info/ftm.list
/usr/sbin/ftm
/lib/compress_vart.sh
/etc/init.d/ftm
