# usr/lib/opkg/info/iptables-mod-fullconenat.list
/usr/lib/iptables/libipt_FULLCONENAT.so
/etc/init.d/xt_fullcone
