# usr/lib/opkg/info/kmod-qca-wifi-unified-profile.list
/ini/wifi_module_param.ini
/ini/QCN9000.ini
/lib/wifi/qcacommands_ol_radio.xml
/ini/internal/QCA5018_i.ini
/ini/QCA5018.ini
/usr/sbin/thermaltool
/usr/sbin/exttool
/ini/internal/QCA6290_i.ini
/usr/sbin/athadhoc
/ini/QCN6122.ini
/usr/lib/libwlanfw_wmi.so
/ini/QCA8074.ini
/lib/modules/4.4.60/ath_pktlog.ko
/lib/modules/4.4.60/rawmode_sim.ko
/etc/config/qcacfg80211
/etc/init.d/ssid_steering
/usr/sbin/cfg80211tool_mesh
/ini/internal/QCA8074V2_i.ini
/ini/internal/QCA9984_i.ini
/usr/sbin/apstats
/lib/modules/4.4.60/qdf.ko
/ini/QCA9984.ini
/usr/sbin/spectraltool
/lib/wifi/tools_config
/usr/sbin/peerratestats
/lib/wifi/qcacommands_mesh.xml
/lib/modules/4.4.60/qca_ol.ko
/ini/QCA6018.ini
/ini/internal/global_i.ini
/lib/wifi/debug/counters_wlan_prompt_Dakota.sh
/lib/modules/4.4.60/wifi_2_0.ko
/etc/config/ssid-steering
/ini/internal/QCA6018_i.ini
/usr/sbin/assocdenialnotify
/usr/sbin/pktlogdump
/usr/sbin/80211stats
/usr/sbin/cfg80211tool
/lib/modules/4.4.60/smart_antenna.ko
/usr/lib/libqca_tools.so
/usr/sbin/qldtool
/ini/internal/IPQ4019_i.ini
/usr/sbin/radartool
/ini/internal/QCA8074_i.ini
/lib/modules/4.4.60/wifi_3_0.ko
/lib/wifi/qcawifi_countrycode.txt
/etc/init.d/boot-ftm
/usr/sbin/wifistats
/usr/sbin/wifitool
/usr/sbin/wifi_hw_mode
/usr/lib/libtlvparser.so
/ini/internal/QCN9000_i.ini
/usr/sbin/acsdbgtool
/usr/lib/libtlvencoder.so
/usr/sbin/athtestcmd
/ini/internal/AR900B_i.ini
/usr/sbin/tx99tool
/lib/modules/4.4.60/asf.ko
/usr/sbin/athstats
/usr/lib/libapstats.so
/usr/lib/libtlvtemplate.so
/lib/wifi/qcawificfg80211.sh
/lib/wifi/debug/counters_wlan_prompt_akronite.sh
/usr/lib/libtlvcmdrsp.so
/usr/lib/libwlanfw_htt.so
/lib/wifi/qca-wifi-modules
/etc/init.d/diag_socket_app
/ini/internal/QCA9888_i.ini
/lib/sbin/wifi_try
/ini/internal/QCN6122_i.ini
/usr/lib/libtlvutil.so
/usr/sbin/wlanconfig
/lib/wifi/qcacommands_vap.xml
/lib/wifi/debug/counters_wlan_prompt_hawkeye.sh
/lib/upgrade/keep.d/wifi_nss
/ini/global.ini
/usr/sbin/cfr_test_app
/usr/sbin/qca_gensock
/etc/hotplug.d/firmware/11-wifi-fw
/usr/lib/libwifistats.so
/lib/modules/4.4.60/umac.ko
/usr/sbin/athstatsclr
/usr/sbin/ssidsteering
/lib/modules/4.4.60/qca_spectral.ko
/usr/sbin/pktlogconf
/usr/sbin/wps_enhc
/ini/QCA6290.ini
/lib/modules/4.4.60/mem_manager.ko
/lib/wifi/qcawifi.sh
/usr/sbin/wifi_try
/usr/sbin/athkey
