# usr/lib/opkg/info/luci-app-arpbind.list
/usr/lib/lua/luci/controller/arpbind.lua
/usr/lib/lua/luci/model/cbi/arpbind.lua
/etc/init.d/arpbind
/etc/uci-defaults/luci-arpbind
/etc/hotplug.d/iface/50-arpbind
