# usr/lib/opkg/info/luci-app-ssr-plus.list
/usr/lib/lua/luci/view/shadowsocksr/reset.htm
/usr/lib/lua/luci/view/shadowsocksr/subscribe.htm
/usr/lib/lua/luci/model/cbi/shadowsocksr/client.lua
/usr/lib/lua/luci/view/shadowsocksr/ssrurl.htm
/usr/share/shadowsocksr/gfw2ipset.sh
/usr/lib/lua/luci/view/shadowsocksr/checkport.htm
/etc/uci-defaults/luci-ssr-plus
/usr/share/shadowsocksr/chinaipset.sh
/etc/ssrplus/gfw_list.conf
/usr/lib/lua/luci/model/cbi/shadowsocksr/advanced.lua
/usr/lib/lua/luci/model/cbi/shadowsocksr/server-config.lua
/usr/lib/lua/luci/view/shadowsocksr/certupload.htm
/usr/share/shadowsocksr/genred2config.sh
/usr/share/shadowsocksr/ssrplusupdate.sh
/usr/bin/ssr-switch
/etc/ssrplus/ad.conf
/usr/lib/lua/luci/model/cbi/shadowsocksr/servers.lua
/etc/ssrplus/china_ssr.txt
/etc/ssrplus/white.list
/usr/lib/lua/luci/model/cbi/shadowsocksr/server.lua
/usr/lib/lua/luci/view/shadowsocksr/refresh.htm
/etc/ssrplus/black.list
/etc/ssrplus/netflixip.list
/etc/init.d/shadowsocksr
/etc/ssrplus/oversea_list.conf
/etc/ssrplus/deny.list
/usr/lib/lua/luci/view/shadowsocksr/ping.htm
/etc/ssrplus/gfw_base.conf
/etc/config/shadowsocksr
/usr/share/rpcd/acl.d/luci-app-ssr-plus.json
/usr/share/shadowsocksr/update.lua
/usr/lib/lua/luci/model/cbi/shadowsocksr/log.lua
/usr/lib/lua/luci/model/cbi/shadowsocksr/status.lua
/lib/upgrade/keep.d/luci-app-ssr-plus
/usr/lib/lua/luci/model/cbi/shadowsocksr/client-config.lua
/usr/share/shadowsocksr/subscribe.lua
/usr/share/shadowsocksr/gen_config.lua
/usr/bin/ssr-monitor
/usr/lib/lua/luci/model/cbi/shadowsocksr/control.lua
/usr/lib/lua/luci/view/shadowsocksr/server_list.htm
/etc/ssrplus/netflix.list
/usr/lib/lua/luci/view/shadowsocksr/status.htm
/usr/lib/lua/luci/view/shadowsocksr/check.htm
/usr/lib/lua/luci/controller/shadowsocksr.lua
/usr/lib/lua/luci/view/shadowsocksr/socket.htm
/usr/lib/lua/luci/view/shadowsocksr/cache.htm
/usr/bin/ssr-rules
