# usr/lib/opkg/info/qca-qrtr.list
/usr/sbin/qrtr-ns
/usr/sbin/qrtr-lookup
/usr/sbin/qrtr-cfg
/etc/init.d/qrtr
