# usr/lib/opkg/info/qca-acfg.list
/etc/acfg_common.conf
/usr/sbin/acfg_tool
/etc/init.d/qca-acfg
