# usr/lib/opkg/info/shadowsocks-libev-config.list
/etc/init.d/shadowsocks-libev
/etc/config/shadowsocks-libev
