# usr/lib/opkg/info/udpxy.list
/usr/bin/udpxrec
/usr/bin/udpxy
/etc/config/udpxy
/etc/init.d/udpxy
