# usr/lib/opkg/info/vsftpd-alt.list
/etc/config/vsftpd
/usr/sbin/vsftpd_wrapper
/usr/sbin/vsftpd_prepare
/usr/sbin/vsftpd
/lib/upgrade/keep.d/vsftpd-alt
/etc/init.d/vsftpd
